# ngix/root/etc/services.d/nginx/nginx.conf
[program:crond]
command=nginx -g "daemon off;"
